# rpiclock.service
[Unit]
Description=rpiclock

[Service]
ExecStart=/home/pi/rpiclock.py
User=pi

[Install]
WantedBy=default.target
