# rpiclock.service
[Unit]
Description=rpiclock

[Service]
ExecStart=/home/pi/rpiclock

[Install]
WantedBy=default.target
