# rpiclock.service
[Unit]
Description=rpiclock

[Service]
ExecStart=/usr/local/bin/rpiclock
User=pi

[Install]
WantedBy=default.target
